# cron/crontab
*  *  *  *  * /bin/sync_upstream_repo sync_repo >> /var/log/sync_repo.log 2>&1
