# cron/crontab
*  *  *  *  * /bin/sync_upstream_repo sync_repo >> /tmp/log.log 2>> /tmp/errors.log
